# deploy/clov-api.service
# Clov API systemd 유닛 — /etc/systemd/system/clov-api.service 에 복사
# WorkingDirectory에 clov-api.jar + application-secret.yaml + clov-api.env 를 함께 둔다.
# (application-secret.yaml은 커밋 금지 시크릿 — VM에서 직접 생성. Spring이 working dir에서 자동 로드)
[Unit]
Description=Clov API (Spring Boot)
After=network.target mysql.service
Wants=mysql.service

[Service]
User=clov
Group=clov
WorkingDirectory=/opt/clov-api
EnvironmentFile=/opt/clov-api/clov-api.env
ExecStart=/usr/bin/java -jar /opt/clov-api/clov-api.jar
# Spring Boot 정상 종료 코드
SuccessExitStatus=143
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
